# service/src/fridgilyse.service
[Unit]
Description=MQTT Fridgilyse Analysis Daemon

[Service]
ExecStart=/root/fridgilyse.py

[Install]
WantedBy=multi-user.target


